# skel/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/local/bin/run.py
StandardInput=tty
StandardOutput=tty
